# testdata/propagation/deadlock/primary.service
[Unit]
Description=Primary service
BindsTo=secondary.service
After=secondary.service

[Service]
ExecStart=/bin/true
