# testdata/propagation/deadlock/secondary.service
[Unit]
Description=Secondary service
Requires=primary.service
After=primary.service

[Service]
ExecStart=/bin/true
