# testdata/graph/after_without_requires/database.service
[Unit]
Description=Database

[Service]
ExecStart=/usr/bin/db
