# testdata/graph/after_without_requires/web.service
[Unit]
Description=Web server
After=database.service

[Service]
ExecStart=/usr/bin/web
