# testdata/graph/dangling_requires/app.service
[Unit]
Description=Application service
Requires=missing-db.service
After=missing-db.service

[Service]
ExecStart=/usr/bin/app
