# testdata/validation/timer_no_trigger/empty.timer
[Unit]
Description=Timer without any triggers

[Timer]
# No OnCalendar, OnBootSec, etc.

[Install]
WantedBy=timers.target
