# testdata/graph/cycle_simple/a.service
[Unit]
Description=Service A
Requires=b.service

[Service]
ExecStart=/bin/true
