# testdata/graph/cycle_simple/b.service
[Unit]
Description=Service B
Requires=c.service

[Service]
ExecStart=/bin/true
