# testdata/graph/cycle_simple/c.service
[Unit]
Description=Service C
Requires=a.service

[Service]
ExecStart=/bin/true
