# testdata/graph/requires_without_after/app.service
[Unit]
Description=Application
Requires=cache.service

[Service]
ExecStart=/usr/bin/app
