# testdata/graph/requires_without_after/cache.service
[Unit]
Description=Cache service

[Service]
ExecStart=/usr/bin/cache
