# testdata/validation/socket_no_service/orphan.socket
[Unit]
Description=Socket without matching service

[Socket]
ListenStream=8080

[Install]
WantedBy=sockets.target
