# testdata/propagation/restart_storm/a.service
[Unit]
Description=Service A
BindsTo=b.service
After=b.service

[Service]
ExecStart=/bin/true
Restart=on-failure
