# testdata/propagation/restart_storm/b.service
[Unit]
Description=Service B
BindsTo=a.service
After=a.service

[Service]
ExecStart=/bin/true
Restart=on-failure
