# testdata/graph/linear_chain/s1.service
[Unit]
Description=Service 1

[Service]
ExecStart=/bin/true
