# testdata/graph/linear_chain/s2.service
[Unit]
Description=Service 2
Requires=s1.service
After=s1.service

[Service]
ExecStart=/bin/true
