# testdata/graph/linear_chain/s3.service
[Unit]
Description=Service 3
Requires=s2.service
After=s2.service

[Service]
ExecStart=/bin/true
