# testdata/validation/service_valid/good.service
[Unit]
Description=A valid service

[Service]
Type=simple
ExecStart=/bin/true
User=root

[Install]
WantedBy=multi-user.target
