# testdata/validation/service_missing_exec/bad.service
[Unit]
Description=Service missing ExecStart

[Service]
Type=simple
User=root

[Install]
WantedBy=multi-user.target
